# docker/supervisor/supervisord.conf
<?xml version="1.0" encoding="UTF-8"?>
<!-- Supervisor configuration for Laravel Queue Workers -->
<!-- Place in /etc/supervisor/conf.d/ or include via /etc/supervisor/conf.d/*.conf -->

[supervisorctl]
serverurl=unix:///var/run/supervord.sock

[unix_http_server]
file=/var/run/supervord.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface_factory

[supervisor:server]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface_factory
supervisor.control.enablefcgi = true

; ============================================================
; Laravel Queue Workers
; ============================================================

; OTP Queue Worker - High Priority
[program:laravel-queue-otp]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work redis --queue=otp,default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/app/storage/logs/supervisor_otp.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
priority=100

; Default Queue Worker
[program:laravel-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work redis --queue=default,emails --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/app/storage/logs/supervisor_default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
priority=200

; Notification Queue Worker
[program:laravel-queue-notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work redis --queue=notifications,websocket --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/app/storage/logs/supervisor_notifications.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
priority=300